# curator4-1974071052.conf
[program:curator4-1974071052]
command=python3 /home/curator4/curator4-1974071052/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1974071052/test.err.log
stdout_logfile=/home/curator4/curator4-1974071052/test.out.log